# conf/supervisor.conf
[program:transmitter]
directory=/home/harold/relay
command=python3 transmitter.py
autorestart=true
startretries=1000000
