# app/batch.service
[Unit]
After=db.service
Requires=db.service

[Service]
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker rm -f batch
ExecStartPre=/usr/bin/docker run -d --name batch --link db app python3 -m app.batch
ExecStart=/usr/bin/docker logs -f batch
ExecStop=/usr/bin/docker stop batch

[Install]
WantedBy=multi-user.target
